# rootfs/etc/services.d/openclaw/run
#!/bin/bash
# OpenClaw gateway service
source /etc/s6-overlay/lib/env-utils.sh
source_env_prefix OPENCLAW_

echo "[openclaw] Starting openclaw gateway..."

# Show generated config
echo "[openclaw] DEBUG: Config file contents:"
if [ -f /data/.openclaw/openclaw.json ]; then
    cat /data/.openclaw/openclaw.json
else
    echo "[openclaw] DEBUG: /data/.openclaw/openclaw.json NOT FOUND"
fi

# Allow openclaw to use tailscale
if [ "${TAILSCALE_ENABLE:-false}" = "true" ]; then
    tailscale set --operator=openclaw
fi

# Background port checker - monitors ALL listening ports
(
    for i in $(seq 1 30); do
        sleep 3
        echo "[openclaw] PORT-CHECK ($((i*3))s): All listening ports:"
        cat /proc/net/tcp 2>/dev/null | awk 'NR>1{split($2,a,":"); port=strtonum("0x"a[2]); if($4=="0A") print "  listening on port " port " (local addr: " a[1] ")"}' 2>/dev/null || echo "  cannot read /proc/net/tcp"
    done
) &

# First, discover what gateway options exist
echo "[openclaw] DEBUG: Running 'openclaw gateway --help':"
runuser -l openclaw -c '
    export NVM_DIR="$HOME/.nvm"
    [ -s "$NVM_DIR/nvm.sh" ] && source "$NVM_DIR/nvm.sh"
    export PNPM_HOME="$HOME/.local/share/pnpm"
    export PATH="$PNPM_HOME:$PATH"
    openclaw gateway --help 2>&1
' 2>&1
echo "[openclaw] DEBUG: --help exit code: $?"

echo "[openclaw] DEBUG: About to start gateway..."

# Run gateway directly as openclaw user
runuser -l openclaw -c '
    source /etc/s6-overlay/lib/env-utils.sh 2>/dev/null
    source_env_prefix OPENCLAW_ 2>/dev/null
    export NVM_DIR="$HOME/.nvm"
    [ -s "$NVM_DIR/nvm.sh" ] && source "$NVM_DIR/nvm.sh"
    export PNPM_HOME="$HOME/.local/share/pnpm"
    export PATH="$PNPM_HOME:$PATH"

    echo "[openclaw] DEBUG: Executing: openclaw gateway --allow-unconfigured"
    echo "[openclaw] DEBUG: OPENCLAW_STATE_DIR=$OPENCLAW_STATE_DIR"
    echo "[openclaw] DEBUG: OPENCLAW_GATEWAY_TOKEN=$OPENCLAW_GATEWAY_TOKEN"

    openclaw gateway --allow-unconfigured 2>&1
    echo "[openclaw] DEBUG: Gateway process exited with code: $?"
' 2>&1

EXIT_CODE=$?
echo "[openclaw] DEBUG: runuser exited with code: $EXIT_CODE"

# If gateway exits, sleep before restart to avoid rapid restart loop
sleep 5
